# systemd/wallpaper.service
[Unit]
Description=Randomly change the wallpaper

[Service]
Type=oneshot
ExecStart=/home/dave/.dotfiles/scripts/randomwallhaven jess
